# etc/udev/rules.d/90-monitor_switch.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/opt/monitor_udev_switch.sh"
